# sdk/rootfs/root.released/usr/etc/inittab
# < id >:< runlevels >:< action >:< process >
::sysinit:/bin/mount -a
::sysinit:/usr/etc/init.d/rcS
::respawn:-/bin/ash

